# Distribution/LaunchAgent/com.joymapkit.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joymapkit.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BINARY_PATH__</string>
        <string>run</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <!-- Log paths use __HOME__ placeholder; substitute with actual home directory at install time -->
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/JoyMapKit/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/JoyMapKit/stderr.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
